# roles/fastly_origin/files/fastly-invalidate@.path
[Unit]
Description=Watching /srv/ftp/%i/os/x86_64/%i.db.tar.gz for changes

[Path]
PathChanged=/srv/ftp/%i/os/x86_64/%i.db.tar.gz
