# roles/gitlab_runner/files/gitlab-runner-docker-cleanup.timer
[Unit]
Description=Run gitlab-runner-docker-cleanup.service daily

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1d

[Install]
WantedBy=timers.target
