# roles/gitlab_runner/files/libvirt-executor-vm-template.service
[Unit]
Description=Create updated VM image with the required tools

[Service]
Type=oneshot
ExecStart=/usr/local/bin/libvirt-executor create-vm-template
