# roles/dbscripts/files/arch-svntogit.service
[Unit]
Description=arch-svntogit

[Service]
Type=oneshot
User=svntogit
ExecStart=/srv/svntogit/update-repos.sh
