# roles/dbscripts/files/arch-svntogit.timer
[Unit]
Description=arch-svntogit

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
