# roles/dbscripts/files/gen_rsyncd.timer
[Unit]
Description=Generate rsync config for mirrors

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
