# roles/dbscripts/files/integrity-check.timer
[Unit]
Description=Integrity Check

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
