# roles/dbscripts/files/lastsync.service
[Unit]
Description=lastsync

[Service]
Type=oneshot
User=ftp
ExecStart=/usr/bin/bash -c "date '+\%s\' > /srv/ftp/lastsync"
