# roles/dbscripts/files/lastsync.timer
[Unit]
Description=lastsync

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
