# roles/dbscripts/files/sourceballs.timer
[Unit]
Description=Sourceballs

[Timer]
OnBootSec=15min
OnUnitActiveSec=8h

[Install]
WantedBy=timers.target
