# roles/dbscripts/files/update-abs.service
[Unit]
Description=update-abs

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/update-abs.sh
