# roles/codesearch/templates/codesearch-index.timer.j2
[Unit]
Description=Codesearch reindex timer

[Timer]
OnUnitInactiveSec=5min
Unit=codesearch-index.service

[Install]
WantedBy=timers.target
