# roles/archweb/templates/archweb-populate_signoffs.timer.j2
[Unit]
Description=archweb populate_signoffs timer

[Timer]
OnUnitActiveSec=20m
OnBootSec=15min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
