# roles/archweb/templates/archweb-rebuilderd.timer.j2
[Unit]
Description=archweb rebuilderd timer

[Timer]
OnUnitActiveSec=60m
OnBootSec=15min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
