# roles/archweb/templates/archweb-mirrorcheck.timer.j2
[Unit]
Description=archweb mirrorcheck

[Timer]
OnUnitActiveSec=1h
OnBootSec=10min
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
