# roles/archweb/templates/archweb-donor_import.timer.j2
[Unit]
Description=archweb donor_import timer

[Timer]
OnUnitActiveSec=1w
OnBootSec=15min
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
