# roles/archweb/templates/archweb-bumpbuddy.timer.j2
[Unit]
Description=archweb bumpbuddy timer

[Timer]
OnUnitActiveSec=60m
OnBootSec=15min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
