# roles/spampd/files/sa-update.timer
[Unit]
Description=sa-update

[Timer]
OnCalendar=*-*-* 00:05:00

[Install]
WantedBy=timers.target
