# roles/spampd/files/sa-update.service
[Unit]
Description=sa-update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sa-update.sh
