# roles/aurweb/templates/aurweb-git.service.j2
[Unit]
Description=Git gc AUR

[Service]
Type=Oneshot
User={{ aurweb_user }}
WorkingDirectory={{ aurweb_git_dir }}
ExecStart=/usr/bin/git gc
Nice=5
