# roles/aurweb/templates/aurweb-votereminder.timer.j2
[Unit]
Description=Send aurweb vote reminder emails

[Timer]
OnStartupSec=12h
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
