# roles/aurweb/templates/aurweb-aurblup.timer.j2
[Unit]
Description=Update aurweb blacklist

[Timer]
OnStartupSec=2h
OnUnitActiveSec=2h
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
