# roles/aurweb/templates/aurweb-github-mirror.timer.j2
[Timer]
OnCalendar=minutely
RandomizedDelaySec=1m

[Install]
WantedBy=timers.target
