# roles/aurweb/templates/aurweb-pkgmaint.timer.j2
[Unit]
Description=Delete old, empty reserved aurweb pkgbases

[Timer]
OnStartupSec=2h
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
