# roles/aurweb/templates/aurweb-git.timer.j2
[Unit]
Description=Git gc & repack AUR

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
