# roles/aurweb/templates/aurweb-tuvotereminder.service.j2
[Unit]
Description=Send aurweb TU vote reminder emails
Requires=mysqld.service
After=mysqld.service

[Service]
Type=oneshot
User={{ aurweb_user }}
ExecStart=/usr/local/bin/aurweb-tuvotereminder
