# roles/aurweb/templates/aurweb-aurblup.service.j2
[Unit]
Description=Update aurweb blacklist
Requires=mysqld.service
After=mysqld.service

[Service]
Type=Oneshot
User={{ aurweb_user }}
ExecStart=/usr/local/bin/aurweb-aurblup
