# roles/aurweb/templates/aurweb-popupdate.service.j2
[Unit]
Description=Update aurweb per-package popularity counts
Requires=mysqld.service
After=mysqld.service

[Service]
Type=Oneshot
User={{ aurweb_user }}
ExecStart=/usr/local/bin/aurweb-popupdate
