# roles/aurweb/templates/aurweb-mkpkglists.service.j2
[Unit]
Description=Regenerate aurweb package and user lists
Requires=mysqld.service
After=mysqld.service

[Service]
Type=Oneshot
User={{ aurweb_user }}
ExecStart=/usr/local/bin/aurweb-mkpkglists
