# roles/aurweb/templates/aurweb-git.service.j2
[Unit]
Description=Git gc AUR

[Service]
Type=oneshot
User={{ aurweb_user }}
WorkingDirectory={{ aurweb_dir }}
ExecStart=/usr/local/bin/aurweb-git-gc.sh
Nice=5
