# roles/docker_image/templates/docker-image.timer.j2
[Unit]
Description=Timer for docker-image.service

[Timer]
OnCalendar={{ docker_image_time }}
Persistent=true
Unit=docker-image.service

[Install]
WantedBy=timers.target
