# roles/aur/templates/aurweb-git.timer.j2
[Unit]
Description=Git gc & repack AUR

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
