# roles/gitlab/files/gitlab-bot-token-extender.timer
[Unit]
Description=GitLab Bot Token Extender

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=24h

[Install]
WantedBy=timers.target
