# roles/borg-client/files/borg-backup.service
[Unit]
Description=Borg backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/borg-backup.sh
