# roles/bumpbuddy/files/bumpbuddy.timer
[Unit]
Description=Schedule for Upstream Release Tracker

[Timer]
OnCalendar=00/03:00
Persistent=true

[Install]
WantedBy=timers.target
