# roles/bumpbuddy/files/repo-syncer.service
[Service]
Type=oneshot
User=bumpbuddy
EnvironmentFile=/etc/conf.d/bumpbuddy
ExecStart=/opt/bumpbuddy/repo-syncer.py
