# roles/dovecot/files/dovecot-cleanup.service
[Unit]
Description=Daily Dovecot cleanup

[Service]
Type=oneshot
ExecStart=/usr/bin/doveadm purge -A
