# roles/archmanweb/templates/archmanweb_update.timer.j2
[Unit]
Description=Timer for the archmanweb update

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
