# roles/kanboard/templates/kanboard-cron.timer
[Unit]
Description=Kanboard cronjob

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
