# roles/prometheus_exporters/templates/prometheus-flyspray-textcollector.timer.j2
[Unit]
Description=Prometheus Flyspray Exporter TextCollector Timer

[Timer]
OnUnitActiveSec=5m
OnBootSec=2min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
