# roles/prometheus_exporters/templates/prometheus-smart-textcollector.timer.j2
[Unit]
Description=Prometheus Smart Exporter TextCollector Timer

[Timer]
OnUnitActiveSec=1h
OnBootSec=15min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
