# roles/prometheus_exporters/templates/prometheus-fail2ban-textcollector.timer.j2
[Unit]
Description=Prometheus Fail2ban Exporter TextCollector Timer

[Timer]
OnUnitActiveSec=15s
OnBootSec=1min

[Install]
WantedBy=timers.target
