# roles/prometheus_exporters/templates/prometheus-aur-textcollector.timer.j2
[Unit]
Description=Prometheus Arch Exporter TextCollector Timer

[Timer]
OnUnitActiveSec=5m
OnBootSec=5min

[Install]
WantedBy=timers.target
