# roles/prometheus_exporters/templates/prometheus-repository-textcollector.timer.j2
[Unit]
Description=Prometheus Arch Repository Exporter TextCollector Timer

[Timer]
# 24 hours
OnUnitActiveSec=1440m
OnBootSec=15min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
