# roles/prometheus_exporters/templates/prometheus-arch-textcollector.timer.j2
[Unit]
Description=Prometheus Arch Exporter TextCollector Timer

[Timer]
OnUnitActiveSec=60m
OnBootSec=15min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
