# roles/prometheus_exporters/templates/prometheus-btrfs-textcollector.timer.j2
[Unit]
Description=Prometheus Btrfs Exporter TextCollector Timer

[Timer]
OnUnitActiveSec=1d
OnBootSec=15min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
