# roles/aurweb/files/aurweb-popupdate.timer
[Unit]
Description=Update aurweb per-package popularity counts

[Timer]
OnStartupSec=2h
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
