# roles/aur/files/aurweb-aurblup.timer
[Unit]
Description=Update aurweb blacklist

[Timer]
OnStartupSec=2h
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
