# roles/borg_client/templates/borg-backup.service.j2
[Unit]
Description=Borg backup

[Service]
Type=oneshot
Restart=on-failure
RestartSec=5m
{% for backup in backup_hosts %}
ExecStart=/usr/local/bin/borg-backup{{ backup['suffix'] }}.sh
{% endfor %}
