# roles/borg_client/templates/borg-backup-offsite.service.j2
[Unit]
Description=Borg backup (offsite)

[Service]
Type=oneshot
{% include 'borg-backup-cpu-quota.j2' %}
ExecStart=/usr/local/bin/borg-backup-offsite.sh
