# roles/prometheus_exporters/files/gitlab-exporter.timer
[Unit]
Description=Gitlab Exporter Timer

[Timer]
OnUnitActiveSec=10min
OnBootSec=5min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
