# roles/archweb/files/archweb-mirrorcheck.timer
[Unit]
Description=archweb mirrorcheck

[Timer]
OnUnitActiveSec=1h
OnBootSec=10min
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
