# roles/db/files/blacklist-purger.service
[Service]
ExecStart=/usr/local/bin/blacklist-purger

[Install]
WantedBy=multi-user.target
