# roles/borg_client/files/borg-backup-offsite.service
[Unit]
Description=Borg backup (offsite)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/borg-backup-offsite.sh
