# roles/arch-boxes/templates/arch-boxes.timer.j2
[Unit]
Description=Timer for arch-boxes.service

[Timer]
OnCalendar={{ archboxes_time }}
Persistent=true
Unit=arch-boxes.service

[Install]
WantedBy=timers.target
