# roles/arch-boxes/templates/arch-boxes.service.j2
[Unit]
Description=arch-boxes service for deploying images and vagrant boxes

[Service]
Type=oneshot
ExecStart=/usr/bin/packer-io
User={{ archboxes_user }}
WorkingDirectory={{ archboxes_dir }}
