# roles/archbuild/files/clean-offload-build.timer
[Unit]
Description=Daily offload-build artifacts cleanup

[Timer]
OnCalendar=daily
RandomizedDelaySec=1d
Persistent=true

[Install]
WantedBy=timers.target
