# roles/archbuild/files/clean-dests.timer
[Unit]
Description=Hourly dest cleanup

[Timer]
OnCalendar=hourly
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
