# roles/archbuild/files/mkpkg@.timer
[Unit]
Description=Nightly build of %I

[Timer]
OnCalendar=22:00
RandomizedDelaySec=8h

[Install]
WantedBy=timers.target
