# roles/fail2ban/files/blacklist-syncer.service
[Service]
ExecStart=/usr/local/bin/blacklist-syncer

[Install]
WantedBy=multi-user.target
