# roles/syncdebug/files/syncdebug.timer
[Unit]
Description=Sync debug packages every minute

[Timer]
OnCalendar=minutely
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
