# roles/arch_boxes/templates/arch-boxes.service.j2
[Unit]
Description=arch-boxes service for deploying images and vagrant boxes

[Service]
Type=oneshot
ExecStart={{ archboxes_git_dir }}/controller.py
User={{ archboxes_user }}
